# deploy/wmux-mediamtx.service.example
[Unit]
Description=wmux MediaMTX WebRTC stream router
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/mediamtx %h/.wmux/mediamtx.yml
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
